# ext/solaris/smf/pxp-agent.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="pxp-agent">

  <service name="network/pxp-agent" type="service" version="1">

    <create_default_instance enabled="false"/>
    <single_instance/>

    <dependency name="loopback" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/network/loopback:default"/>
    </dependency>

    <dependency name="physical" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/network/physical:default"/>
    </dependency>

    <dependency name="fs-local" grouping="require_all" restart_on="none" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <exec_method type="method" name="start" exec="/opt/puppetlabs/puppet/bin/pxp-agent" timeout_seconds="60"/>

    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>

    <exec_method type="method" name="refresh" exec=":kill" timeout_seconds="60"/>

    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">PCP Execution Protocol (PXP) Agent</loctext>
      </common_name>
      <documentation>
        <manpage title="pxp-agent" section="1"/>
        <doc_link name="puppetlabs.com" uri="http://puppetlabs.com/puppet/introduction"/>
      </documentation>
    </template>
  </service>

</service_bundle>
